# etc/systemd/system/kvmd-oled-reboot.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/kvmd-oled --height=32 --offset-y=6 --interval=0 --text="Rebooting...\nPlease wait"
